# data/bamfdaemon.service.in
[Unit]
Description=BAMF Application Matcher Framework
PartOf=graphical-session.target

[Service]
Type=dbus
BusName=org.ayatana.bamf
ExecStart=@PKGLIBEXECDIR@/bamfdaemon
Restart=on-failure
